# deploy/systemd/diary-backup.timer.example
[Unit]
Description=Run Telegram Diary SQLite Backup Daily

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=diary-backup.service

[Install]
WantedBy=timers.target
